# root_overlay/etc/cron.d/osm-backup
SHELL=/bin/bash
43 02 * * * root /opt/devtank/backup.sh >> /var/log/osm-backup
